# DEBIAN/postinst
#!/bin/bash

if [  ! -f /etc/wghrs/config.toml ]
then
   cp /etc/wghrs/default-config.toml /etc/wghrs/config.toml
fi

systemctl enable wghrs
